# scripts/controlpanel-owlcms.service
# systemd unit — primary OWLCMS instance with Tracker
#
# Starts both OWLCMS and Tracker headlessly using the default (main) instance.
# Uses the same restart semantics as Docker's on-failure policy:
#   Exit 0       → do not restart
#   Exit non-zero → restart (used intentionally on database loads)
#   Deliberate external stop (TERM/KILL) → do not restart
#   Abnormal signal (SIGSEGV, …) → restart
#
# Installation:
#   cp controlpanel-owlcms.service ~/.config/systemd/user/
#   systemctl --user daemon-reload
#   systemctl --user enable --now controlpanel-owlcms.service
#
# Logs:
#   journalctl --user -u controlpanel-owlcms.service -f

[Unit]
Description=OWLCMS Competition Management (main instance)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
# --- adjust User to the account that owns ~/.local/share/owlcms-controlpanel ---
User=jflamy

# Launch both OWLCMS and Tracker with the last-run (or latest installed) versions.
# Add explicit version numbers if needed, e.g. --owlcms 65.0.0 --tracker 2.3.0
# we do not use --mqtt on this instance
ExecStart=/usr/bin/controlpanel --owlcms --tracker

# Stop: SIGTERM → Java shutdown hooks → clean exit.
ExecStop=/usr/bin/controlpanel --owlcms stop --tracker stop
KillMode=mixed
TimeoutStopSec=15

# Restart policy matching Docker on-failure:
Restart=on-failure
SuccessExitStatus=SIGTERM
RestartSec=3

# Harden
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
# The control panel needs write access to its data directory
ReadWritePaths=%h/.local/share/owlcms-controlpanel

[Install]
WantedBy=default.target